# site/bin/site-app-supervisor.conf
[program:cms-site-app]
stdout_logfile=/var/log/cms/site.log
stderr_logfile=/var/log/cms/site.log
nodaemon=false
user = cmsapp
directory = /usr/share/cms/site/bin
command=perl -w app.pl
