# admin/bin/admin-app-supervisor.conf
[program:cms-admin-app]
stdout_logfile=/var/log/cms/admin.log
stderr_logfile=/var/log/cms/admin.log
nodaemon=false
user = cmsapp
directory = /usr/share/cms/admin/bin
command=perl -w app.pl
